# docker/root/etc/services.d/squid-service/run
#!/usr/bin/with-contenv bash
s6-setuidgid abc java -jar /app/squid-service/squid-service.jar server /config/config.yml
